# etc/supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log

[supervisorctl]

[inet_http_server]
port=127.0.0.1:9001
username=admin
password=admin

[program:test]
command=python 
autostart=false
